# files/mysql_supervisor.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mysql]
# user=nobody
command=/usr/bin/mysqld_safe
